# deploy/systemd/o2o@.service
[Unit]
Description=O2O Carpooling %i
After=network-online.target docker.service
Wants=network-online.target docker.service

[Service]
Type=simple
WorkingDirectory=/var/www/o2o-Local-Carpooling
EnvironmentFile=/var/www/o2o-Local-Carpooling/.env
Environment="SPRING_PROFILES_ACTIVE=demo"
Environment="NACOS_DISCOVERY_ENABLED=false"
Environment="SPRING_CLOUD_NACOS_DISCOVERY_ENABLED=false"
Environment="SPRING_CLOUD_NACOS_DISCOVERY_REGISTER_ENABLED=false"
Environment="SPRING_CLOUD_SERVICE_REGISTRY_AUTO_REGISTRATION_ENABLED=false"
Environment="SERVER_ADDRESS=127.0.0.1"
Environment="SPRING_CLOUD_NACOS_DISCOVERY_IP=127.0.0.1"
Environment="MANAGEMENT_HEALTH_SENTINEL_ENABLED=false"
Environment="SPRING_CLOUD_SENTINEL_ENABLED=false"
Environment="MANAGEMENT_HEALTH_DISCOVERY_ENABLED=false"
Environment="MANAGEMENT_HEALTH_DEFAULTS_ENABLED=false"
Environment="MANAGEMENT_HEALTH_DB_ENABLED=false"
Environment="MANAGEMENT_HEALTH_REDIS_ENABLED=false"
Environment="MANAGEMENT_HEALTH_RABBIT_ENABLED=false"
Environment="MANAGEMENT_HEALTH_MONGO_ENABLED=false"
Environment="MANAGEMENT_HEALTH_MAIL_ENABLED=false"
Environment="SPRING_CLOUD_OPENFEIGN_CLIENT_CONFIG_DEFAULT_CONNECT_TIMEOUT=10000"
Environment="SPRING_CLOUD_OPENFEIGN_CLIENT_CONFIG_DEFAULT_READ_TIMEOUT=90000"
Environment="SPRING_CLOUD_INETUTILS_TIMEOUT_SECONDS=1"
Environment="MYSQL_JDBC_URL=jdbc:mysql://127.0.0.1:3307/o2o_carpooling?useUnicode=true&characterEncoding=utf8&serverTimezone=Asia/Shanghai"
Environment="HIKARI_MAXIMUM_POOL_SIZE=2"
Environment="HIKARI_MINIMUM_IDLE=0"
Environment="HIKARI_IDLE_TIMEOUT=30000"
Environment="HIKARI_CONNECTION_TIMEOUT=120000"
EnvironmentFile=-/var/www/o2o-Local-Carpooling/deploy/systemd/env/%i.env
ExecStart=/usr/bin/java -Xms24m -Xmx96m -XX:MaxMetaspaceSize=96m -XX:ReservedCodeCacheSize=24m -Xss512k -XX:+UseSerialGC -XX:ActiveProcessorCount=1 -XX:TieredStopAtLevel=1 -Dspring.jmx.enabled=false -Dfile.encoding=UTF-8 -jar /var/www/o2o-Local-Carpooling/backend/%i/target/%i-0.5.0-SNAPSHOT.jar --spring.datasource.url=jdbc:mysql://127.0.0.1:3307/o2o_carpooling?useUnicode=true&characterEncoding=utf8&serverTimezone=Asia/Shanghai --spring.datasource.username=${MYSQL_USERNAME} --spring.datasource.password=${MYSQL_PASSWORD} --spring.datasource.hikari.maximum-pool-size=${HIKARI_MAXIMUM_POOL_SIZE} --spring.datasource.hikari.minimum-idle=${HIKARI_MINIMUM_IDLE} --spring.datasource.hikari.idle-timeout=${HIKARI_IDLE_TIMEOUT} --spring.datasource.hikari.connection-timeout=${HIKARI_CONNECTION_TIMEOUT}
Restart=on-failure
RestartSec=10
SuccessExitStatus=143
NoNewPrivileges=true
PrivateTmp=true
MemoryAccounting=true
MemoryHigh=180M
MemoryMax=260M
TasksMax=256

[Install]
WantedBy=multi-user.target
